# com.herald.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.herald</string>

    <key>ProgramArguments</key>
    <array>
        <string>__JAVA_BIN__</string>
        <string>-jar</string>
        <string>__HOME__/.herald/herald-bot.jar</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/.herald</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/herald.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/herald.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HERALD_TELEGRAM_BOT_TOKEN</key>
        <string>__HERALD_TELEGRAM_BOT_TOKEN__</string>
        <key>ANTHROPIC_API_KEY</key>
        <string>__ANTHROPIC_API_KEY__</string>
        <key>HERALD_TELEGRAM_ALLOWED_CHAT_ID</key>
        <string>__HERALD_TELEGRAM_ALLOWED_CHAT_ID__</string>
        <key>OPENAI_API_KEY</key>
        <string>__OPENAI_API_KEY__</string>
        <key>GCAL_MCP_URL</key>
        <string>__GCAL_MCP_URL__</string>
        <key>GMAIL_MCP_URL</key>
        <string>__GMAIL_MCP_URL__</string>
    </dict>
</dict>
</plist>
